# utils/dashboard/gunicorn.conf
[program:gunicorn]
command=gunicorn dashboard:app -b localhost:5000
directory=/deploy/app